# roles/valkey/files/valkey-server.service
[Unit]
Description=Valkey data structure server
Documentation=https://github.com/valkey-io/valkey-doc
Wants=network-online.target
After=valkey-kerberos-ticket-refresher.service

[Service]
ExecStart=/opt/valkey/bin/valkey-server /opt/valkey/etc/valkey.conf
LimitNOFILE=10032
NoNewPrivileges=yes
Type=notify
TimeoutStartSec=infinity
TimeoutStopSec=infinity
UMask=0077
User=valkey
Group=valkey

[Install]
WantedBy=multi-user.target