# roles/caddy/files/caddy-restarter.path
[Path]
Unit=caddy-restarter.service
PathChanged=/config/services.d

[Install]
WantedBy=multi-user.target